# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

; 单进程跑 admin(8001)+proxy(8999)，避免双份 Python 与重复模块常驻（RSS 主要来源）
[program:app]
command=python src/main.py both
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
